# cron.d/getepg
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

29 */2 * * *   www-data	/var/www/epgrec/getepg.php
